# axon/bin/AxonIvyEngine.service
[Unit]
Description=Axon.ivy Engine
Wants=network.target
After=network.target

[Service]
User=<user>
Group=<group>
Type=forking
PIDFile=/<installationdir>/bin/AxonIvyEngine.pid
WorkingDirectory=/<installationdir>
ExecStart=/<installationdir>/bin/AxonIvyEngineService.sh startdaemon
ExecStop=/<installationdir>/bin/AxonIvyEngineService.sh stop
#RestartSec=60
#Restart=on-failure
StartLimitInterval=0
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
